# deploy/systemd/ohlcv-cleanup.timer
[Unit]
Description=Investment Advisor - OHLCV cleanup 주간 실행 (매주 일요일 08:00 KST)

[Timer]
# 매주 일요일 08:00 KST — 06:30 메인 배치 + 07:30 meta-sync 완료 후 retention 정리
OnCalendar=Sun *-*-* 08:00:00 Asia/Seoul
Persistent=true
Unit=ohlcv-cleanup.service

[Install]
WantedBy=timers.target
