# deploy/systemd/pre-market-briefing.service
[Unit]
Description=Investment Advisor - 프리마켓 브리핑 (KST 06:30, 미국 장 마감 직후)
Documentation=https://github.com/yawnsduzin/investment-advisor
After=network-online.target postgresql.service
Wants=network-online.target
# 가격/OHLCV sync가 먼저 끝난 뒤 브리핑 — 같은 06:30 트랜잭션에서 순서 보장
After=universe-sync-price.service
After=universe-sync-indices.service

[Service]
Type=oneshot
User=__SYSTEM_USER__
Group=__SYSTEM_USER__
WorkingDirectory=__INSTALL_DIR__
EnvironmentFile=__INSTALL_DIR__/.env
Environment=PYTHONUNBUFFERED=1
Environment=PATH=__INSTALL_DIR__/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=__VENV_PYTHON__ -m analyzer.briefing_main
# Claude SDK + DB 집계 — 30분 여유
TimeoutStartSec=1800
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=5

Restart=no
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
