# deploy/systemd/universe-sync-meta.timer
[Unit]
Description=Investment Advisor - Universe 메타 주간 실행 (매주 일요일 03:30 KST)

[Timer]
OnCalendar=Sun *-*-* 03:30:00 Asia/Seoul
Persistent=true
Unit=universe-sync-meta.service

[Install]
WantedBy=timers.target
