# deploy/systemd/universe-sync-price.timer
[Unit]
Description=Investment Advisor - Universe 가격/OHLCV 매일 06:30 KST 실행

[Timer]
# 매일 06:30 KST — 미국 장 마감 직후(EDT 16:00→KST 05:00 / EST 16:00→KST 06:00)
# KRX+US OHLCV 일괄 확보. pre-market-briefing.service /
# investment-advisor-analyzer.service 가 `After=universe-sync-price.service`
# 의존성으로 이 unit 완료 후 자동 실행됨.
OnCalendar=*-*-* 06:30:00 Asia/Seoul
Persistent=true
Unit=universe-sync-price.service

[Install]
WantedBy=timers.target
